# recipes-core/goldiservices/files/goldi-rauc-updater.service
[Unit]
Description = The GOLDi RAUC Updater Script

[Service]
ExecStart = goldi-rauc-updater.sh

[Install]
WantedBy = multi-user.target
